# rootfs/etc/services.d/varnish/run
#!/usr/bin/with-contenv sh

exec service varnish start

#exec systemctl restart varnish.service